# systemd/ews-skill-sync.service
[Unit]
Description=EWS Skill Stdio JSON-RPC Daemon
After=network-online.target
Wants=network-online.target

[Service]
# User/Group are injected by scripts/install.sh during install.
Type=simple
WorkingDirectory=/opt/ews-skill
EnvironmentFile=/opt/ews-skill/.env
RuntimeDirectory=ews-skill
RuntimeDirectoryMode=0755
ExecStart=/opt/ews-skill/ews_skilld --transport unix --socket /run/ews-skill/daemon.sock
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
